# nomadnet-pages/server-config/nomadnet.service
[Unit]
Description=NomadNet — ZeroPoint Presence Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nomadnet --daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
